# roles/desktop/templates/systemd/restic_backup.timer.j2
[Unit]
Description=Run Restic Backup
After=network-online.target
[Timer]
OnCalendar=*-*-* 21:30:00
Persistent=true
[Install]
WantedBy=timers.target
