# roles/ambition/templates/iperf3.service.j2
[Unit]
Description=iperf3
After=network-online.target

[Service]
User={{main_username}}
ExecStart=/usr/bin/iperf3 -s -f M

Restart=always
RestartSec=15

[Install]
WantedBy=multi-user.target